# crontab
0 5  * * *  /home/ubuntu/microsvc-letsencrypt-docker/cron_job.sh
